# venus/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user command
*/20  * *   *   *   root . /etc/planet_env; update-page 2>&1 | logger

